# supervisord.conf
[unix_http_server]
file=/var/run/supervisord/supervisord.sock

[supervisord]
pidfile=/var/run/supervisord/supervisord.pid
logfile=/var/log/supervisord.log
loglevel=info
user=root
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord/supervisord.sock

[program:mysqld]
command=/usr/local/bin/mariadb_run.sh
user=mysql
autostart=true

[program:nginx]
command=/usr/sbin/nginx
user=root
autostart=true

[program:php-fpm]
command=/usr/sbin/php-fpm7 -F
user=root
autostart=true
